# deploy/gateway.service
[Unit]
Description=noVPN legacy compatibility wrapper for vpn-service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/novpn/vpn-service
ExecStart=/usr/local/bin/vpn-service -config /etc/novpn/vpn-service/config.yaml
Restart=always
RestartSec=2
LimitNOFILE=200000

[Install]
WantedBy=multi-user.target
